# files/etc/crontabs/root
# ============================================================
# D525 路由器 — 定时任务
# ============================================================
# 格式: 分 时 日 月 周 命令

# 每周日凌晨 4:00 更新 OpenAppFilter 特征库
0 4 * * 0 /usr/bin/oaf-update

# 每天凌晨 3:30 重启 MWAN3（防止连接表泄漏）
30 3 * * * /etc/init.d/mwan3 restart >/dev/null 2>&1
